# debian/qvpnapplet/etc/xdg/autostart/qvpnapplet.desktop
[Desktop Entry]
Type=Application
Exec=/usr/bin/qvpnapplet
Name=qvpnapplet

