# debian/qvpnapplet.service
[Unit]
Description=qVPNd
After=network.target

[Service]
ExecStart=/usr/sbin/qvpnd

[Install]
WantedBy=multi-user.target
